# setup/systemd/oar-server.service.in
[Unit]
Description=OAR server
Documentation=man:oar-server(1)
After=network-online.target
After=remote-fs.target
After=postgresql.service
After=mysql-server.service
Wants=network-online.target

[Service]
ExecStart=%%SBINDIR%%/oar-server
Type=forking

[Install]
WantedBy=multi-user.target
