# setup/systemd/oar-node.service.in
[Unit]
Description=OAR compute node
After=network-online.target
After=remote-fs.target
Wants=network-online.target

[Service]
ExecStart=/usr/sbin/sshd -f %%OARCONFDIR%%/sshd_config
Type=notify
KillMode=process
RuntimeDirectory=sshd
RuntimeDirectoryMode=0755
RuntimeDirectoryPreserve=yes

[Install]
WantedBy=oar.target
