# setup/systemd/oar.target
[Unit]
Description=OAR target
After=multi-user.target
AllowIsolate=yes

[Install]
RequiredBy=multi-user.target
