# setup/cron.d/oar-server.in
#
# cron.d/oar-server -- schedules OAR accounting data mining
# 
# $Id: oar-server.cron.d 1235 2008-03-18 10:03:12Z neyron $

# By default, run every hours
0 * * * * %%ROOTUSER%% [ -x %%SBINDIR%%/oaraccounting ] && %%SBINDIR%%/oaraccounting
