# setup/cron.d/oar-node.in
#
# cron.d/oar-node -- schedules periodic checks of the local node by running
# every scripts of the /etc/oar/check.d directory, if any.
# 
# $Id: oar-node.cron.d 1235 2008-03-18 10:03:12Z neyron $

# By default, run every hours
0 * * * * root [ -x %%OARDIR%%/oarnodecheckrun ] && %%OARDIR%%/oarnodecheckrun
