# debian.notusedanymore/oar-desktop-computing-cgi.cron.hourly
#!/bin/bash
# OAR stagein cache directory clean-up

/usr/sbin/oarcache
