# debian/oar-server.cron.daily
#!/bin/bash

# Run OAR Accounting
/usr/sbin/oaraccounting
