# setup/cron.hourly/oar-desktop-computing-cgi.in
#!/bin/bash
# OAR stagein cache directory clean-up

%%SBINDIR%%/oarcache
