# setup/init.d/oar-desktop-computing-agent.in
#! /bin/sh
#
# desktop-computing-agent
#
# --/--/----: Based on startup scripts from Bruno Bzeznik <bruno.bzeznik@imag.fr>
# 02/09/2011: Merge RPM and Debian scripts by Philippe Le Brouster <philippe.le-brouster@imag.fr>
#

# chkconfig: 2345 90 10
# description: This script starts or stops the desktop computing agent
# processname: oar-agent-daemon
# pidfile: %%OARHOMEDIR%%/oar_agent.pid


### BEGIN INIT INFO
# Provides:         oar-desktop-computing-agent
# Required-Start:   $network $local_fs $remote_fs
# Required-Stop:    $network $local_fs $remote_fs
# Default-Start:    2 3 4 5
# Default-Stop:     0 1 6
# Short-Description:    OAR desktop computing agent daemon starting
### END INIT INFO

PATH=%%SBINDIR%%:%%BINDIR%%:/sbin:/bin:/usr/sbin:/usr/bin
NAME=oar-desktop-computing-agent
DESC=desktop-computing-agent
DAEMON=%%SBINDIR%%/oar-agent-daemon
DAEMON_NAME=oar-agent-daemon
PIDFILE=%%OARHOMEDIR%%/oar_agent.pid
SCRIPTNAME=%%INITDIR%%/$NAME
NOLSB=

[ -f /lib/lsb/init-functions ] || NOLSB=yes

if [ -f /etc/debian_version ]; then
    system=debian
elif [ -f /etc/redhat-release ]; then
    system=redhat
elif [ -f /etc/SuSE-release ]; then
    system=suse
elif [ -f /etc/gentoo-release ]; then
    system=gentoo
fi

# Exit if the package is not installed
[ -x "$DAEMON" ] || exit 0

if [ -z "$NOLSB" ]; then
    . /lib/lsb/init-functions
    fail_msg() {
        echo ""
        log_failure_msg "$@"
    }
    warn_msg() {
        log_warning_msg "$@"
    }
    succ_msg() {
        log_success_msg "$@"
    }
    begin_msg() {
        echo -n "$@:"
    }
else
    echo "This system doesn't provide the LSB functions. Failing"
    exit 2
fi


do_start() {
    begin_msg "Starting $DESC"
    if pidofproc -p $PIDFILE > /dev/null; then
        fail_msg " already running"
        exit 1
    fi
    start_daemon -p "$PIDFILE" -n "-20" "$DAEMON" start
    RET=$?
    if [ "$RET" != 0 ]; then
        fail_msg "Unable to start"
    else
        succ_msg
    fi
}

do_stop() {
    begin_msg "Stopping $DESC"
    killproc -p $PIDFILE
    RETVAL="$?"
    if [ "$RETVAL" = 2 ]; then
        fail_msg "Unable to stop $DESC"
        exit 2
    fi
    killproc $DAEMON_NAME
    if [ "$?" = 2 ]; then
        fail_msg "Unable to stop $DESC"
        exit 2
    fi
    rm -f $PIDFILE

    succ_msg
    exit "$RETVAL"
}

case "$1" in
    start)
        do_start
        ;;
    stop)
        do_stop
        ;;
    reload|restart|force-reload)
        if do_stop; then
            sleep 1
            do_start
        fi
        ;;
    *)
        echo "Usage: $SCRIPTNAME {start|stop|restart|force-reload}" >&2
        exit 1
        ;;
esac


exit 0
